# linux/.config/autostart/compton.desktop
[Desktop Entry]
Name=Compton
Exec=compton
Type=Application
Terminal=false
X-GNOME-Autostart-enabled=true
